# services/hyperion-router.service
[Unit]
Description=Hyperion Router - Telegram to Claude Code bridge
After=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/hyperion
Environment=PATH=/home/admin/.local/bin:/home/admin/.cargo/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/admin/hyperion/config/config.env
ExecStart=/home/admin/hyperion/.venv/bin/python /home/admin/hyperion/src/bot/hyperion_bot.py
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyperion-router

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
