# services/hyperion-daemon.service
[Unit]
Description=Hyperion Daemon - Always-on Claude Code message processor
After=network.target hyperion-router.service
Wants=hyperion-router.service

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/hyperion-workspace
Environment=PATH=/home/admin/.claude/bin:/home/admin/.local/bin:/home/admin/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/admin
ExecStart=/usr/bin/python3 /home/admin/hyperion/src/daemon/daemon.py
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyperion-daemon

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
